# contrib/systemd/strata-pvr-scheduler.timer
[Unit]
Description=Run Strata PVR scheduler periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
AccuracySec=1min
Persistent=true
Unit=strata-pvr-scheduler.service

[Install]
WantedBy=timers.target
